# deploy/table-tennis-dashboard.service
[Unit]
Description=Table Tennis Forecast dashboard (HTTP server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/table-tennis-forecast
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/root/table-tennis-forecast/.env
ExecStart=/root/table-tennis-forecast/.venv/bin/python /root/table-tennis-forecast/app/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
